# packages/mentor-event-bus/plists/com.caia.mentor.memory-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.caia.mentor.memory-watcher
  ==============================
  Recursive fs.watch on agent/memory/. Emits MemoryWritten on every
  detected create / modify / delete (debounced at 500ms per path).

  Installed by: packages/mentor-event-bus/scripts/install.sh
  Logs:         ~/Library/Logs/caia-mentor/memory-watcher.{out,err}.log

  Substitutions handled by install.sh:
    __USER_HOME__       — $HOME
    __REPO__            — absolute path to caia checkout
    __EVENTS_DB_PATH__  — events.sqlite path
    __MEMORY_DIR__      — directory to watch (default $REPO/agent/memory)
    __NODE_BIN__        — node binary path
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.caia.mentor.memory-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>__NODE_BIN__</string>
        <string>__REPO__/packages/mentor-event-bus/dist/cli.js</string>
        <string>watch-memory</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__USER_HOME__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>CAIA_EVENT_BUS_DB_PATH</key>
        <string>__EVENTS_DB_PATH__</string>
        <key>CAIA_MEMORY_DIR</key>
        <string>__MEMORY_DIR__</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__USER_HOME__/Library/Logs/caia-mentor/memory-watcher.out.log</string>

    <key>StandardErrorPath</key>
    <string>__USER_HOME__/Library/Logs/caia-mentor/memory-watcher.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
