# apps/executor/com.conductor.executor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Conductor Executor Daemon launchd plist
  Install: launchctl load ~/Library/LaunchAgents/com.conductor.executor.plist
  Stop:    launchctl unload ~/Library/LaunchAgents/com.conductor.executor.plist
  Logs:    ~/Documents/conductor-logs/executor.log

  KeepAlive: true — daemon restarts automatically on crash.
  RunAtLoad: false — must be enabled explicitly via `conductor exec start`.
  ThrottleInterval: 10 — minimum 10s between restarts to prevent crash loops.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.conductor.executor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>CONDUCTOR_CLI_PATH</string>
    <string>exec</string>
    <string>daemon</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>CONDUCTOR_API</key>
    <string>http://localhost:7776</string>
  </dict>

  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>CONDUCTOR_LOG_DIR/executor.log</string>

  <key>StandardErrorPath</key>
  <string>CONDUCTOR_LOG_DIR/executor.err.log</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
