# packages/apprentice-corpus/plists/com.chiefaia.apprentice-corpus.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  com.chiefaia.apprentice-corpus.plist

  LaunchAgent that runs the Apprentice corpus aggregator daily at
  02:00 local time. Outputs go to ~/Documents/projects/apprentice/corpora/<YYYY-MM-DD>/.

  Install via the package's install script (renders + bootstraps):
    pnpm --filter @chiefaia/apprentice-corpus build
    packages/apprentice-corpus/scripts/install-apprentice-corpus.sh <orchestrator-session-id>

  The script substitutes CURRENT_NODE_BIN, CURRENT_PKG_DIR, CURRENT_HOME and
  CURRENT_SESSION_ID before bootstrapping launchd.

  Uninstall:
    launchctl bootout gui/$(id -u)/com.chiefaia.apprentice-corpus
    rm ~/Library/LaunchAgents/com.chiefaia.apprentice-corpus.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chiefaia.apprentice-corpus</string>

  <key>ProgramArguments</key>
  <array>
    <string>CURRENT_NODE_BIN</string>
    <string>CURRENT_PKG_DIR/dist/cli.js</string>
    <string>aggregate</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- Patched at install time by scripts/install-apprentice-corpus.sh -->
    <key>CAIA_MEMORY_DIR</key>
    <string>CURRENT_HOME/Library/Application Support/Claude/local-agent-mode-sessions/CURRENT_SESSION_ID/agent/memory</string>

    <key>CAIA_REPORTS_DIR</key>
    <string>CURRENT_HOME/Documents/projects/reports</string>

    <key>CAIA_EVENTS_DB</key>
    <string>CURRENT_HOME/.caia/mentor/events.sqlite</string>

    <key>CAIA_GITHUB_REPO</key>
    <string>prakashgbid/caia</string>

    <key>APPRENTICE_CORPUS_ROOT</key>
    <string>CURRENT_HOME/Documents/projects/apprentice/corpora</string>

    <key>CLAUDE_BINARY_PATH</key>
    <string>CURRENT_CLAUDE_BIN</string>

    <key>PATH</key>
    <string>CURRENT_PATH</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>CURRENT_HOME/Library/Logs/chiefaia/apprentice-corpus.log</string>

  <key>StandardErrorPath</key>
  <string>CURRENT_HOME/Library/Logs/chiefaia/apprentice-corpus.err.log</string>

  <key>RunAtLoad</key>
  <false/>

  <key>WorkingDirectory</key>
  <string>CURRENT_PKG_DIR</string>

  <key>ProcessType</key>
  <string>Background</string>

  <key>Nice</key>
  <integer>5</integer>
</dict>
</plist>
