# packages/ai-first-upgrade-scanner/launchd/com.caia.ai-first-upgrade-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.caia.ai-first-upgrade-daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/opt/node@22/bin/node</string>
        <string>/Users/macbook32/Documents/projects/caia/packages/ai-first-upgrade-scanner/bin/ai-first-upgrade-scan</string>
        <string>--quiet</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/macbook32/Library/Logs/caia/ai-first-upgrade.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/macbook32/Library/Logs/caia/ai-first-upgrade.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/macbook32</string>
    </dict>
</dict>
</plist>
