# packages/secrets-broker/com.stolution.secrets-broker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Launchd plist for the secrets broker HTTP service on stolution (macOS).

  Install:
    cp com.stolution.secrets-broker.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.stolution.secrets-broker.plist

  On Linux (stolution server), use the equivalent systemd unit instead.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.stolution.secrets-broker</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/home/s903/projects/plugins/secrets-broker/dist/server.cjs</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>BROKER_PORT</key>
    <string>7788</string>
    <key>BROKER_HOST</key>
    <string>127.0.0.1</string>
    <key>BROKER_LOG_LEVEL</key>
    <string>info</string>
    <key>VAULT_ADAPTER</key>
    <string>ssh-file</string>
    <key>VAULT_SSH_HOST</key>
    <string>stolution</string>
    <key>VAULT_DIR</key>
    <string>/home/s903/.vault</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/secrets-broker.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/secrets-broker-error.log</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
