# apps/local-preview-orchestrator/plists/com.stolution.local-preview.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.stolution.local-preview.dashboard
  =====================================
  Site supervisor for the CAIA dashboard preview.

  Reads the live build via the `current` symlink at:
    __USER_HOME__/Library/Application Support/Stolution/local-preview/dashboard/current/

  Listens on http://localhost:5173.

  KeepAlive=true so the process is restarted by launchd whenever the deploy
  daemon SIGTERMs it after a symlink swap.

  Wait-for-current shim (PR-F): on first boot, the `current` symlink does
  not exist until the deploy daemon has run a successful build. Rather than
  exiting non-zero (which crash-loops via KeepAlive=Crashed=true on a 10s
  ThrottleInterval), we busy-poll for the symlink with `until ... do sleep 5`
  and only `cd` into it once it exists. The wait is bounded only by the
  agent eventually being booted out, which is the right shape for a
  long-lived supervisor.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stolution.local-preview.dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>INSTALL_DIR="$LOCAL_PREVIEW_INSTALL_ROOT/dashboard"; CURRENT="$INSTALL_DIR/current"; echo "[supervisor:dashboard] waiting for $CURRENT" &gt;&amp;2; until [ -L "$CURRENT" ]; do sleep 5; done; echo "[supervisor:dashboard] symlink found; starting" &gt;&amp;2; cd "$CURRENT" || exit 1; echo $$ &gt; "$INSTALL_DIR/pid"; exec pnpm --filter @caia-app/dashboard exec next start -p 5173</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__USER_HOME__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>LOCAL_PREVIEW_INSTALL_ROOT</key>
        <string>__USER_HOME__/Library/Application Support/Stolution/local-preview</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__USER_HOME__/Library/Logs/local-preview/dashboard.out.log</string>

    <key>StandardErrorPath</key>
    <string>__USER_HOME__/Library/Logs/local-preview/dashboard.err.log</string>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
