# apps/local-preview-orchestrator/plists/com.stolution.local-preview.deploy-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.stolution.local-preview.deploy-daemon
  =========================================
  Long-running poll daemon. Iterates every 30s through all configured sites
  and invokes deploySite() for any whose origin/<branch> SHA differs from
  the live `current` symlink.

  Installed by: apps/local-preview-orchestrator/scripts/install.sh
  Source:       apps/local-preview-orchestrator/plists/com.stolution.local-preview.deploy-daemon.plist
  Logs:         ~/Library/Logs/local-preview/deploy-daemon.{out,err}.log

  Substitute `__USER_HOME__` with $HOME and `__REPO__` with the absolute
  path to your caia checkout when installing. The install.sh script does
  this automatically.

  Per-site branch overrides: __LOCAL_PREVIEW_<SITE>_BRANCH__ placeholders
  are sed-substituted at install time from the operator's env (or the
  bake-in defaults of develop/master/main respectively).

  Mentor event-bus integration (PR-H, leg-4 stage-6 finding):
  CAIA_EVENT_BUS_DB_PATH points the orchestrator's lazily-opened
  MentorClient at the production events.sqlite managed by the
  mentor-event-bus install. If Mentor is not installed yet, the client
  open silently no-ops; the orchestrator retries opening on each
  successful deploy so it picks up Mentor whenever it eventually arrives.
  See apps/local-preview-orchestrator/src/cli.ts `tryOpenMentorClient`.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stolution.local-preview.deploy-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>__REPO__/apps/local-preview-orchestrator/dist/cli.js</string>
        <string>poll-loop</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__USER_HOME__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>LOCAL_PREVIEW_INSTALL_ROOT</key>
        <string>__USER_HOME__/Library/Application Support/Stolution/local-preview</string>
        <key>LOCAL_PREVIEW_BUILD_WORKSPACE</key>
        <string>/private/tmp/local-preview-build</string>
        <key>LOCAL_PREVIEW_DASHBOARD_BRANCH</key>
        <string>__LOCAL_PREVIEW_DASHBOARD_BRANCH__</string>
        <key>LOCAL_PREVIEW_POKER_ZENO_BRANCH</key>
        <string>__LOCAL_PREVIEW_POKER_ZENO_BRANCH__</string>
        <key>LOCAL_PREVIEW_ROULETTE_COMMUNITY_BRANCH</key>
        <string>__LOCAL_PREVIEW_ROULETTE_COMMUNITY_BRANCH__</string>
        <key>CAIA_EVENT_BUS_DB_PATH</key>
        <string>__CAIA_EVENT_BUS_DB_PATH__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__USER_HOME__/Library/Logs/local-preview/deploy-daemon.out.log</string>

    <key>StandardErrorPath</key>
    <string>__USER_HOME__/Library/Logs/local-preview/deploy-daemon.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
